# data/systemd/shopsystem-mail-daily-invoice.service.in
[Unit]
Description=Send shopsystem daily invoice mail
DefaultDependencies=no
Conflicts=shutdown.target

[Service]
Type=oneshot
User=shop
Group=shop
ExecStart=@PATH@/shop-invoice temporary
StandardOutput=journal
StandardError=journal
IOSchedulingClass=idle
