# data/systemd/shopsystem-ui.service.in
[Unit]
Description=Shop System UI

[Service]
User=shop
Group=shop
Type=simple
ExecStart=@PATH@/shop-curses-ui
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty3
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

[Install]
WantedBy=multi-user.target
