# data/systemd/shopsystem-mail-backup.service.in
[Unit]
Description=Send shopsystem backup mail
DefaultDependencies=no
Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=oneshot
User=shop
Group=shop
ExecStart=@PATH@/shop-backup
IOSchedulingClass=idle
