# data/systemd/shopsystem-display-on-off.service.in
[Unit]
Description=Display on-off via MQTT
After=systemd-user-sessions.service

[Service]
Type=simple
ExecStart=@PATH@/display-on-off.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
