# data/systemd/shopsystem-input-device.service.in
[Unit]
Description=Shop System Input Device Handler
BindsTo=dev-barcodescanner.device
Requires=dev-barcodescanner.device
After=dev-barcodescanner.device

[Service]
Type=dbus
BusName=io.mainframe.shopsystem.InputDevice
User=shop
Group=shop
Restart=always
RestartSec=10
ExecStart=@PATH@/shop-input-device
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty3

[Install]
WantedBy=multi-user.target
