# systemd/shopsystem-report.timer
[Unit]
Description=Send shopsystem daily invoice

[Timer]
OnCalendar=*-*-* 08:10:00
Persistent=true

[Install]
WantedBy=timers.target
